# etc/local.d/v4l2.start
#!/bin/bash

# We shouldn't need this, the /etc/modules or whatever should just load it, but I'll put it here just in case

modprobe v4l2loopback
